# retired/launchagents/com.id8labs.friday-prep-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.id8labs.friday-prep-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eddiebelaval/.hydra/tools/friday-prep-sync.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/eddiebelaval/Development/id8/engagements/rose-brill/FRIDAY-PREP.html</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/eddiebelaval/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/eddiebelaval</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/eddiebelaval/Library/Logs/friday-prep-sync/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/eddiebelaval/Library/Logs/friday-prep-sync/launchd-error.log</string>
</dict>
</plist>
